# packages/ddnss/ddnss-1.0.5.tar.gz/ddnss-1.0.5/ddnss@.timer
[Unit]
Description=Hourly update ddnss.de entry for %i

[Timer]
OnCalendar=hourly
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
