# packages/ddnss/ddnss-1.0.5.tar.gz/ddnss-1.0.5/ddnss@.service
[Unit]
Description=Update ddnss.de entry for %i
Requires=network-online.target
After=network-online.target

[Service]
User=ddnss
Group=ddnss
ExecStart=/usr/bin/ddnssupd %i
